# sv/tomcat/run
#!/bin/bash

sv start guacd || exit 1

source /etc/envvars

exec 2>&1
exec /tomcat/bin/catalina.sh run
